# bmc/mmc/obmc-flash-mmc-setprimary@.service.in
[Unit]
Description=Set %I as primary partition

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/usr/bin/obmc-flash-bmc mmc-setprimary %i
