# bmc/mmc/obmc-flash-mmc@.service.in
[Unit]
Description=Write image %I to BMC storage

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/usr/bin/obmc-flash-bmc mmc %i @IMG_UPLOAD_DIR@
