# bmc/mmc/obmc-flash-mmc-mirroruboot.service.in
[Unit]
Description=Copy uboot from the currently booted bmc chip to the alternate chip

[Service]
Type=oneshot
RemainAfterExit=no
Environment=MODE=@STORAGE_MODE@
ExecStart=/usr/bin/obmc-flash-bmc mmc-mirroruboot
