# obmc-flash-host-bios@.service.in
[Unit]
Description=Flash Host Bios image %I to Host

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=obmc-flash-bios vrom ImageHost.pnor %i
