# reboot-guard-disable.service.in
[Unit]
Description=Removes the guard that blocks BMC reboot
After=reboot-guard-enable.service

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/usr/bin/obmc-flash-bmc rebootguarddisable
