# static/force-reboot.service
[Unit]
Description=Force reboot of the BMC

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/sbin/reboot
