# static/obmc-flash-bmc-alt@.service.in
[Unit]
Description=Flash image-bmc to the alt chip

[Service]
Type=oneshot
RemainAfterExit=no
ExecStartPre=-/usr/bin/obmc-flash-bmc umount-static-altfs alt-rofs
ExecStartPre=-/usr/bin/obmc-flash-bmc umount-static-altfs alt-rwfs
ExecStart=/usr/sbin/flashcp /tmp/images/%i/image-bmc /dev/mtd/alt-bmc
ExecStartPost=-/usr/bin/obmc-flash-bmc static-altfs squashfs alt-rofs
ExecStartPost=-/usr/bin/obmc-flash-bmc static-altfs jffs2 alt-rwfs
