# static/obmc-flash-bmc-alt@.service.in
[Unit]
Description=Flash image-bmc to the alt chip and reset cs0

[Service]
Type=oneshot
RemainAfterExit=no
ExecStartPre=-/bin/systemctl stop xyz.openbmc_project.Software.Sync.service
ExecStartPre=-/usr/bin/obmc-flash-bmc umount-static-altfs rofs-alt
ExecStartPre=-/usr/bin/obmc-flash-bmc umount-static-altfs rwfs-alt
ExecStart=/usr/sbin/flashcp /tmp/images/%i/image-bmc /dev/mtd/alt-bmc
ExecStartPost=-/usr/bin/obmc-flash-bmc static-altfs squashfs alt-rofs rofs-alt
ExecStartPost=-/usr/bin/obmc-flash-bmc static-altfs jffs2 alt-rwfs rwfs-alt
ExecStartPost=-/bin/systemctl start xyz.openbmc_project.Software.Sync.service
ExecStartPost=-/usr/bin/reset-cs0-aspeed
