# eeprom-device/xyz.openbmc_project.Software.EEPROMDevice.service.in
[Unit]
Description=Phosphor Firmware Update EEPROM Device

[Service]
ExecStart=/usr/bin/phosphor-fw-update-eeprom-device
Restart=always
Type=dbus
BusName={BUSNAME}

[Install]
WantedBy=multi-user.target
