# bmc/xyz.openbmc_project.Software.Download.service.in
[Unit]
Description=Phosphor Download Manager

[Service]
ExecStart=/usr/libexec/phosphor-bmc-code-mgmt/phosphor-download-manager
Restart=always
Type=dbus
BusName={BUSNAME}

[Install]
WantedBy=multi-user.target
