# bmc/xyz.openbmc_project.Software.Manager.service.in
[Unit]
Description=Phosphor Software Manager

[Service]
ExecStartPre=/bin/sh -c 'mkdir -p /tmp/images'
ExecStart=/usr/bin/phosphor-software-manager
Restart=always
Type=dbus
BusName={BUSNAME}

[Install]
WantedBy=multi-user.target
